# Part76_Cloud_Bot/supervisor-bot.conf
[program:trading-bot]
directory=/root/bot
command=python3 bot.py
user=root
autostart=true
autorestart=false
stderr_logfile=/root/trading-bot.log
stdout_logfile=/root/trading-bot.log